# workspace/crontab/laradock
# * * * * * laradock php /var/www/artisan schedule:run >> /dev/null 2>&1

* * * * * laradock php /var/www/zevo/artisan schedule:run >> /dev/null 2>&1
